# .config/polybar/95-usb.rules
KERNEL=="sd*", ACTION=="add", ATTR{removable}=="1", \
    RUN+="/home/mattia/.config/polybar/system-usb-udev.sh --update"
KERNEL=="sd*", ACTION=="remove", \
    RUN+="/home/mattia/.config/polybar/system-usb-udev.sh --update"
